# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:chromadb]
command=chroma run --host 127.0.0.1 --port 8000 --path /app/data/chroma
autostart=true
autorestart=true
stdout_logfile=/var/log/chromadb.log
stderr_logfile=/var/log/chromadb-error.log
priority=1

[program:web]
command=node packages/web/node_modules/.bin/next start -p 3000
directory=/app
environment=NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/var/log/web.log
stderr_logfile=/var/log/web-error.log
priority=10
startsecs=5

[program:chroma-ui]
command=node chroma-web-ui/node_modules/.bin/next start -p 3001
directory=/app
environment=NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/var/log/chroma-ui.log
stderr_logfile=/var/log/chroma-ui-error.log
priority=11

[program:mcp]
command=node packages/mcp-bailian-rag/dist/http-server.js
directory=/app
environment=NODE_ENV="production",PORT="3002"
autostart=true
autorestart=true
stdout_logfile=/var/log/mcp.log
stderr_logfile=/var/log/mcp-error.log
priority=12

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log
priority=20